# data/systemd/statink-stats-entire-event3.timer
[Unit]
Description=Run stat.ink updater (timer)

[Timer]
OnCalendar=*:12/15
Persistent=false
RandomizedDelaySec=15

[Install]
WantedBy=timers.target
