# data/systemd/statink-stats-update-xpower-distrib3.timer
[Unit]
Description=Run stat.ink updater (timer)

[Timer]
OnCalendar=9:00
Persistent=true
RandomizedDelaySec=2700

[Install]
WantedBy=timers.target
