# data/systemd/statink-gc-cache.timer
[Unit]
Description=Run stat.ink cache gc (timer)

[Timer]
OnCalendar=hourly
RandomizedDelaySec=180

[Install]
WantedBy=timers.target
