# data/systemd/statink-cleanup.timer
[Unit]
Description=Run stat.ink cleanup (timer)

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
