# data/systemd/statink-stats-update-kd-win-rate3.timer
[Unit]
Description=Run stat.ink updater (timer)

[Timer]
OnCalendar=11:00
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
