# data/systemd/statink-dlstats3.timer
[Unit]
Description=Run stat.ink updater (timer)

[Timer]
OnCalendar=9:07
Persistent=true
RandomizedDelaySec=7200

[Install]
WantedBy=timers.target
