# data/systemd/statink-stats-update-splatfest3-power.timer
[Unit]
Description=Run stat.ink updater (timer)

[Timer]
OnCalendar=hourly
Persistent=true
RandomizedDelaySec=50m

[Install]
WantedBy=timers.target
