# data/systemd/statink-fetch-blog.timer
[Unit]
Description=Run stat.ink blog-feed/crawl (timer)

[Timer]
OnCalendar=*:0,15,30,45

[Install]
WantedBy=timers.target
