# data/systemd/statink-update-s2ink.timer
[Unit]
Description=Run stat.ink splatoon2-ink/update (timer)

[Timer]
OnCalendar=*:04

[Install]
WantedBy=timers.target
