# data/systemd/statink-stats-update-entire-salmon2.timer
[Unit]
Description=Run stat.ink stat/update-entire-salmon2 (timer)

[Timer]
OnCalendar=*-*-* 03..21/6:42:00

[Install]
WantedBy=timers.target
