# data/systemd/statink-update-geoip.timer
[Unit]
Description=Run stat.ink geoip/update (timer)

[Timer]
OnCalendar=Thu 03:31

[Install]
WantedBy=timers.target
