# data/systemd/statink-schedule.timer
[Unit]
Description=Run stat.ink scheduling tasks (timer)

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
