# data/systemd/statink-update-s3ink.timer
[Unit]
Description=Run stat.ink splatoon3-ink/update (timer)

[Timer]
OnCalendar=*:02

[Install]
WantedBy=timers.target
