# data/systemd/statink-user-json-auto-update.timer
[Unit]
Description=Run stat.ink user-json/auto-update (timer)

[Timer]
# every minutes
OnCalendar=*-*-* *:*:30

[Install]
WantedBy=timers.target
