# docker/supervisor/h2o.ini
[program:h2o]
command=/usr/bin/h2o -m master -c /etc/h2o/h2o.conf
autostart=true
autorestart=true
