# docker/supervisor/php-fpm.ini
[program:php-fpm]
command=/opt/rh/rh-php56/root/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
