# docker/supervisor/postgres.ini
[program:postgresql]
command=/opt/rh/rh-postgresql94/root/usr/bin/postgres -D /var/opt/rh/rh-postgresql94/lib/pgsql/data
user=postgres
autorestart=true
